# tools/dashboard.service
# Drone dashboard as a user-level systemd service.
#
# Install:
#   mkdir -p ~/.config/systemd/user
#   cp ~/pixhawk/tools/dashboard.service ~/.config/systemd/user/
#   cp ~/pixhawk/tools/dashboard.env.example ~/.config/dashboard.env
#   $EDITOR ~/.config/dashboard.env       # set DASHBOARD_USER / DASHBOARD_PASS
#   systemctl --user daemon-reload
#   systemctl --user enable --now dashboard
#   loginctl enable-linger "$USER"         # so it survives logout/reboot
#
# Logs:    journalctl --user -u dashboard -f
# Restart: systemctl --user restart dashboard

[Unit]
Description=Pixhawk Drone Dashboard (FastAPI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/pixhawk
EnvironmentFile=%h/.config/dashboard.env
# DRONE_CONN, DASHBOARD_USER, DASHBOARD_PASS, DASHBOARD_SECRET, DRONE_CAM,
# DASHBOARD_HOST, DASHBOARD_PORT are all read from EnvironmentFile.
ExecStart=%h/pixhawk/.venv/bin/python -m dashboard.server
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
